# crontab
0 1 * * 2 root /metrics/cron_task collect_metrics.py
