# crontab
@reboot   root /metrics/cron_task /metrics/bin/check_environment.py
0 1 * * 2 root /metrics/cron_task /metrics/bin/collect_metrics.py
