# crontab
0 1 * * 5 root /metrics/process_metrics >> /var/log/cron.log 2>&1
